# debian/marathon-consul.upstart
description "Marathon-consul service (performs Marathon Tasks registration as Consul Services for service discovery)"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
   echo $$ > /var/run/marathon-consul.pid
   exec marathon-consul -config-file=/etc/marathon-consul.d/config.json
end script

post-stop script
    rm -f /var/run/marathon-consul.pid
end script
